# yarra-valley-rice/systemd-units/.config/systemd/user/green-tunnel.service
[Unit]
Description=Green Tunnel service
After=multi-user.target

[Service]
Type=simple
ExecStart=%h/.node_modules/bin/gt --verbose --dns-type "tls" --dns-server "45.90.28.0#8622a7.nextdns.io"
ExecStop=/usr/bin/kill $MAINPID
TimeoutStopSec=5s

[Install]
WantedBy=default.target
