# yarra-valley-rice/systemd-units/.config/systemd/user/xlaunched.target
[Unit]
Description=Custom Target for apps launch
Requires=default.target
After=default.target
AllowIsolate=yes